# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/ANPR_VERCEL.sock

[Install]
WantedBy=sockets.target